# cmd/agent/vpn-server.service
[Unit]
Description=High Performance QUIC VPN Server
After=network-online.target
Wants=network-online.target

[Service]
User=root
Group=root

WorkingDirectory=/etc/vpn-server
ExecStart=/usr/local/bin/vpn-server

Restart=always
RestartSec=3

LimitNOFILE=1048576

# ตั้งค่า Environment Variable (ถ้ามี)
# Environment=JWT_SECRET=your_secret_key

# (Optional) คำสั่งรันก่อนเริ่มแอป (เช่น เคลียร์ Firewall เก่า)
# ExecStartPre=/usr/sbin/iptables -t nat -F

[Install]
WantedBy=multi-user.target